# crontab/date_notify.sh
#!/bin/sh

export DISPLAY=:0.0
date=`date`
notify-send "${date}" -t 2
